# mariadb/files/mariadb.service
.include /usr/lib/systemd/system/mariadb.service

[Service]
LimitNOFILE=infinity
